# open_screen.service
[Unit]
Description=Open Screen
After=multi-user.target[Service]
Type=simple
Restart=always
ExecStart=OpenScreen
WantedBy=multi-user.target